# tcfiles/tc-certmanager.service
[Unit]
Description=Thin Client CertManager
After=NetworkManager-wait-online.service
Wants=NetworkManager-wait-online.service

[Service]
Type=oneshot
RemainAfterExit=yes

# 1. Update Debian CA certificates
ExecStart=/usr/sbin/update-ca-certificates

# 2. Copy custom certs (safe even if /boot/ssl does not exist)
ExecStartPost=/bin/sh -c 'shopt -s nullglob; cp -f /boot/ssl/* /etc/ssl/certs/ 2>/dev/null || true'

# 3. Sync to Citrix cacerts (always safe)
ExecStartPost=/bin/sh -c 'rsync -a --delete /etc/ssl/certs/ /opt/Citrix/ICAClient/keystore/cacerts/'

[Install]
WantedBy=multi-user.target
